# docker/supervisord-ngircd.conf
[program:ngircd]
command = /usr/sbin/ngircd -nf /etc/ngircd/ngircd.conf
user = irc
autostart = true
stderr_logfile = /var/log/ngircd.err
stdout_logfile = /var/log/ngircd.out
